# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:cco]
command=/opt/docker/bin/cco-backend

[program:nginx]
command=/usr/sbin/nginx
